# root/etc/services.d/flood/run
#!/usr/bin/with-contenv bash
#
_term() {
  echo "Caught SIGTERM signal!"
  killall -TERM rtorrent 2>/dev/null
}

trap _term SIGTERM

cd /app/flood
su-exec abc:abc /usr/bin/npm start &

wait
